# supervisor/conf.d/django.conf
[program:django]
directory=server
command=/bin/bash run-server.sh
autostart=true
autorestart=false
startsecs = 0
startretries = 1
redirect_stderr=true
