# supervisor/conf.d/django.conf
[program:django]
directory=server
command=/bin/bash run-server.sh
autostart=true
autorestart=false
redirect_stderr=true
stopasgroup=true
