# supervisor/conf.d/django.conf
[program:django]
directory=server
command=/bin/bash -c ". venv/bin/activate && python manage.py runserver 0.0.0.0:8001"
environment=PATH="venv"
autostart=true
autorestart=false
