# supervisor/conf.d/nginx.conf
[program:nginx]
directory=.
command=nginx -c nginx.conf -p .
autostart=true
autorestart=false
redirect_stderr=true
