# supervisor/conf.d/webpack.conf
[program:webpack]
directory=client
command=yarn start
autostart=true
autorestart=false
