# supervisor/conf.d/yarn.conf
[program:yarn-install]
directory=client
command=yarn install
autostart=false
autorestart=false
redirect_stderr=true
